# deploy/tsx-watcher.service
[Unit]
Description=TSX/TSXV small-cap press release watcher
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Adjust User + paths if you cloned somewhere other than /home/ubuntu/tsx-watcher
User=ubuntu
WorkingDirectory=/home/ubuntu/tsx-watcher
EnvironmentFile=/home/ubuntu/tsx-watcher/deploy/tsx-watcher.env
ExecStart=/usr/bin/env bash /home/ubuntu/tsx-watcher/deploy/run.sh
# The watcher exits cleanly (code 0) at the 4 PM ET close — that is NOT a
# failure, so do not restart it. The timer relaunches it next weekday morning.
# Only restart if it crashes mid-session.
Restart=on-failure
RestartSec=30

# No [Install] section on purpose: this service is launched by the timer, not
# enabled to start on boot (which would re-trigger the 4 PM close repeatedly).
