# deploy/tsx-watcher.timer
[Unit]
Description=Launch the TSX watcher each weekday pre-market (and Sat discovery)

[Timer]
# 06:55 America/New_York, Mon–Sat. systemd handles US Eastern DST automatically.
# Mon–Fri the watcher runs 7:00 AM pre-market → 4:00 PM close in one process.
# Sat it runs the weekly discovery scan (see run.sh).
OnCalendar=Mon-Sat 06:55 America/New_York
# If the VM was off at trigger time, run once on next boot (so a reboot during
# the night doesn't silently skip the day).
Persistent=true
Unit=tsx-watcher.service

[Install]
WantedBy=timers.target
